# meta-scobc-v1/recipes-dualboot/mark-boot-success/files/mark-boot-success.service
[Unit]
Description=Mark successful Linux boot by restoring PMC_MULTI_BOOT register
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/sbin/mark-boot-success
StandardOutput=journal
StandardError=journal
