# meta-scobc-v1/recipes-dualboot/mark-boot-success/files/mark-boot-success.timer
[Unit]
Description=Run mark-boot-success 1min after Linux boot

[Timer]
OnBootSec=1min
Persistent=true
Unit=mark-boot-success.service

[Install]
WantedBy=timers.target
